# Schedule/temp.root.crontab
*/5 * * * * /home/admin/smart-house/Sensor.BTooth/start.cmd
